# deployment/systemd/saas_worker.service
[Unit]
Description=multi_tenant_saas_subscription_platform Worker
After=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/multi_tenant_saas_subscription_platform
EnvironmentFile=/srv/multi_tenant_saas_subscription_platform/.env
ExecStart=/srv/multi_tenant_saas_subscription_platform/.venv/bin/celery -A worker.main:celery_app worker --loglevel=INFO --concurrency=2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
