# deployment/systemd/saas_scheduler.service
[Unit]
Description=multi_tenant_saas_subscription_platform Scheduler
After=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/multi_tenant_saas_subscription_platform
EnvironmentFile=/srv/multi_tenant_saas_subscription_platform/.env
ExecStart=/srv/multi_tenant_saas_subscription_platform/.venv/bin/celery -A worker.scheduler:celery_app beat --loglevel=INFO
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
